# config/crontab
BASE_DIR=/home/pi/AlarmClock

# Crontab to fetch new code every minute
* * * * *   cd BASE_DIR && make install
